# SPECS/consul/consul.service
[Unit]
Description=Consul Agent
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/consul agent -config-dir=/etc/consul.d -dev -advertise 127.0.0.1
ExecReload=/usr/bin/kill -HUP $MAINPID
KillSignal=SIGINT
User=consul
Group=consul
Restart=on-failure
LimitNOFILE=131072

[Install]
WantedBy=multi-user.target
