# SPECS/openssh/sshdat.service
[Unit]
Description=SSH Per-Connection Server

[Service]
ExecStart=/usr/sbin/sshd -i
StandardInput=socket
StandardError=syslog
