# support/cloud-image-builder/etcd.service
[Unit]
Description=etcd service

[Service]
ExecStart=/usr/sbin/etcd
Restart=always

[Install]
WantedBy=multi-user.target
